# deployment_scripts/puppet/modules/lma_infrastructure_alerting_sfdc/templates/sfdc_sender.conf
description "SFDC sender"
author "Max MAzur <sirmax123@gmail.conf>"

start on runlevel [2345]
stop on runlevel [!2345]

chdir /var/run

respawn
respawn limit 20 5
limit nofile 65535 65535

pre-start script
    for i in lock run log lib ; do
	mkdir -p /var/$i/sfdc
	chown root /var/$i/sfdc
    done
end script

script
    [ -x "/usr/lib/nagios/plugins/sfdc/sfdc_sender.py" ] || exit 0
    exec start-stop-daemon --start --chdir /usr/lib/nagios/plugins/sfdc \
	--chuid root:root --make-pidfile --pidfile /var/run/sfdc.pid \
	--exec /usr/lib/nagios/plugins/sfdc/sfdc_sender.py -- -c /usr/lib/nagios/plugins/sfdc/sfdc_nagios.yaml
end script